# templates/etc/profile.d/kafka.sh.j2
export PATH=$PATH:{{kafka_install_parent_dir}}/kafka/bin
