# templates/etc/profile.d/kafka.sh
export PATH=$PATH:/opt/kafka/bin
export KAFKA_ZK_CONNECT="{{kafka_zk_connect}}"
export KAFKA_BROKER_CONNECT="{{kafka_broker_connect}}"
